# dists/reframe@.socket.in
[Unit]
Description=ReFrame Remote Desktop Socket for %i

[Socket]
Accept=no
Service=reframe-streamer@%i.service
SocketUser=@username@
SocketMode=0660
# `DirectoryMode` could not be 0750 because it is owned by systemd as root.
DirectoryMode=0755
ListenStream=%t/reframe/%i.sock

[Install]
WantedBy=sockets.target
