# dists/reframe-server@.service.in
[Unit]
Description=ReFrame Remote Desktop Server for %i
After=network.target
Requires=reframe@%i.socket

[Service]
Type=exec
User=@username@
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
# Required for framebuffer-only rendering.
Environment=EGL_PLATFORM=surfaceless
ExecStart=@bindir@/reframe-server --config=@confdir@/%i.conf --socket=%t/reframe-%i.sock
Restart=on-failure

[Install]
Also=reframe@%i.socket
WantedBy=multi-user.target
